# timer_and_service_files/device_stage_indicator.timer
[Unit]
Description=My Python Script Timer

[Timer]
OnCalendar=*-*-* 08:30:00
Unit=device_stage_indicator.service

[Install]
WantedBy=timers.target