# xaviar_setup/tensorrt_demos/timer_and_service_files/send_count.timer
[Unit]
Description=My Python Script Timer

[Timer]
OnCalendar=*-*-* 08:30:00
Unit=send_count.service

[Install]
WantedBy=timers.target